# public/uploads/folders/support[]4e-503335-20240315113746/report/ericsson-failover-daemon.service
[Unit]
Description=ericsson-failover-daemon
After=network.target

[Service]
Type=simple
User=ericsson
Group=ericsson
# BEGINNING OF BZ47235
PermissionsStartOnly=true
ExecStartPre=/usr/bin/mkdir -p /var/log/ericsson/failover-api/private
ExecStartPre=/usr/bin/chown ericsson: /var/log/ericsson/failover-api/private
# END OF BZ47235
Environment=PYTHONDONTWRITEBYTECODE="True"
ExecStart=/bin/bash -c     "cd /opt/ericsson/failover-api/bin;source activate;./ericsson-failover-daemon --logging /opt/ericsson/failover-api/etc/logging_failover-daemon.json --conf /opt/ericsson/failover-api/etc/failover-api.ini >> /var/log/ericsson/failover-api/private/console_daemon.log 2>&1"
ExecStartPost=/bin/bash -c "ps aux|egrep -e './ericsson-failover-daemon'|grep -v grep|awk '{print $2}'>/var/opt/ericsson/failover-api/ericsson-failover-daemon.pid"
ExecReload=/bin/bash -c    "/bin/pkill -15 -F /var/opt/ericsson/failover-api/ericsson-failover-daemon.pid; sleep 1"
ExecReload=/bin/bash -c    "cd /opt/ericsson/failover-api/bin;source activate;./ericsson-failover-daemon --logging /opt/ericsson/failover-api/etc/logging_failover-daemon.json --conf /opt/ericsson/failover-api/etc/failover-api.ini >> /var/log/ericsson/failover-api/private/console_daemon.log 2>&1;"
ExecReload=/bin/bash -c    "ps aux|egrep -e './ericsson-failover-daemon'|grep -v grep|awk '{print $2}'>/var/opt/ericsson/failover-api/ericsson-failover-daemon.pid"
ExecStop=/bin/bash -c      "/bin/pkill -15 -F /var/opt/ericsson/failover-api/ericsson-failover-daemon.pid; sleep 1"
ExecStopPost=/bin/bash -c  "rm -f /var/opt/ericsson/failover-api/ericsson-failover-daemon.pid"
Restart=always

[Install]
WantedBy=multi-user.target
